# packages/tendril-monitor-vcs/tendril-monitor-vcs-0.1.0.tar.gz/tendril-monitor-vcs-0.1.0/deploy/tendril_vcs_monitor.service

[Unit]
Description=Tendril VCS Monitoring Server

[Service]
ExecStart=/usr/bin/twistd --nodaemon --pidfile= vcs_monitor.tac

WorkingDirectory=/home/tendril/tools/vcs_monitor/deploy

User=tendril
Group=tendril

Restart=always

[Install]
WantedBy=multi-user.target

